# home/private_dot_config/systemd/user/sunset.service.tmpl
[Unit]
Description=Adjust display color temperature by time of day.

[Service]
Type=simple
StandardOutput=journal
ExecStart={{- .chezmoi.homeDir -}}/.config/hypr/scripts/sunset.sh
